# LegendaryOS/Scripts/Kernel-Updater/etc/systemd/system/update-tkg-kernel.timer
[Unit]
Description=LegendaryOS TKG Kernel Updater Timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
